# hbpConfig/~(home)/.config/autostart/working.desktop
[Desktop Entry]
Name=No Name
Type=Application
OnlyShowIn=GNOME;Unity;LXDE
Exec=/home/haobp/.config/autostart/workmnt.sh

